# .config/systemd/user/waybar.service

[Service]
ExecStart=waybar

[Install]
WantedBy=default.target

