# homedir/.config/systemd/user/cliphist.service
[Service]
ExecStart=wl-paste --watch cliphist store

[Install]
WantedBy=default.target

