# homedir/.config/systemd/user/cliphist.service
[Service]
ExecStart=wl-paste --watch /var/home/gmanka/.local/bin/cliphist store

