# homedir/.config/systemd/user/fexserver.service
[Service]
ExecStart=FEXServer -f -p
Restart=always

