# homedir/.config/systemd/user/waybar.service
[Service]
ExecStart=waybar
RestartSec=1s
Restart=always

[Install]
WantedBy=default.target

