# homedir/.config/systemd/user/waybar.service
[Service]
ExecStart=waybar
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target

