# homedir/.config/systemd/user/xwayland-satellite.service.d/override.conf
[Service]
Type=simple
TimeoutStartSec=300

