# homedir/.config/containers/systemd/ygg.container
[Container]
Image=ghcr.io/grisha765/yggdrasil-go
ContainerName=yggdrasil

[Service]
Restart=always

[Install]
WantedBy=default.target

